# 0_SafronovNA/web_nn_service/server/web_nn_server.service
[Unit]
Description=web_nn_server

[Service]
Type=forking

ExecStart=
ExecStop=
ExecReload=
TimeoutSec=300

[Install]
WantedBy=multi-user.target